# server/docker-files/supervisord.conf
[supervisord]
nodaemon=true

[program:powercontrol]
command=/bin/bash -c "cd /app && exec /app/PowerController"
